# .config/autostart/gitboard.desktop
[Desktop Entry]
Type=Application
Name=GitBoard
Exec=/usr/bin/firefox -kiosk https://gitboard.io/dashboard

